# vnc/supervisor/chromium.conf
[program:start-chromium]
command=/usr/local/bin/start-chromium.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/chromium.out.log
stderr_logfile=/var/log/chromium.err.log
priority=60
