# vnc/supervisor/xvfb.conf
[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x800x24
autostart=true
autorestart=true
stdout_logfile=/var/log/xvfb.out.log
stderr_logfile=/var/log/xvfb.err.log
priority=30
